# crontab
*/30 * * * * /home/danny/apps-status-check/app.js
4 10 * * * /home/danny/apps-status-check/app.js daily
